# deploy/pai-worker.service
[Unit]
Description=PAI 主動式 AI 平台 — Queue Worker（認知運行/事件處理）
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=__PHP__ __APP_DIR__/artisan queue:work --sleep=1 --tries=2 --timeout=900
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
